# container-files/etc/profile.d/custom.sh
### Colors
export black="\[\033[0;38;5;0m\]"
export red="\[\033[0;38;5;1m\]"
export orange="\[\033[0;38;5;130m\]"
export green="\[\033[0;38;5;2m\]"
export yellow="\[\033[0;38;5;3m\]"
export blue="\[\033[0;38;5;4m\]"
export bblue="\[\033[0;38;5;12m\]"
export magenta="\[\033[0;38;5;55m\]"
export cyan="\[\033[0;38;5;6m\]"
export white="\[\033[0;38;5;7m\]"
export coldblue="\[\033[0;38;5;33m\]"
export smoothblue="\[\033[0;38;5;111m\]"
export iceblue="\[\033[0;38;5;45m\]"
export turqoise="\[\033[0;38;5;50m\]"
export smoothgreen="\[\033[0;38;5;42m\]"
export defaultcolor="\[\e[m\]"


if [[ ${EUID} == 0 ]] ; then
     PS1="$white[$red\u$coldblue@$green\h$red: $coldblue\w$white]$coldblue\$ $defaultcolor"
else
     PS1="$white[$smoothblue\u$coldblue@$green\h$red: $coldblue\w$white]$coldblue\$ $defaultcolor"
fi

